# rootfs/etc/inittab
::sysinit:/init
tty1::respawn:/sbin/getty -8 38400 tty1
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
